# container-files/etc/supervisord/caddy.conf
[program:caddy]
command=/usr/local/bin/caddy reverse-proxy --from :8282 --to 127.0.0.1:7860
autorestart=true
startsecs=15