# container-files/etc/supervisord/webui.conf
[program:webui]
directory=/home/webui
command=/usr/bin/python3 launch.py --xformers
autorestart=false
redirect_stderr=true
stdout_logfile=/dev/stdout
startsecs=5
user=webui
